# wikipedia/logger/wikipedia-logger.ini
[program:wikipedia-logger-irc]
directory=/1/crawling/wikipedia-logger
command=/1/crawling/wikipedia-logger/start-logger.sh wiki-irc

[program:wikipedia-logger-links]
directory=/1/crawling/wikipedia-logger
command=/1/crawling/wikipedia-logger/start-logger.sh wiki-links
